# core/relay/gipny-relay.service
[Unit]
Description=gipny relay server (i2p)
# Needs the go-i2p SAM router up first (SAMv3 on 127.0.0.1:7656).
After=network-online.target gipny-i2p-router.service
Wants=network-online.target
Requires=gipny-i2p-router.service

[Service]
Type=simple
User=gipny
Group=gipny
WorkingDirectory=/var/lib/gipny-relay
Environment=GIPNY_RELAY_DATA=/var/lib/gipny-relay
ExecStart=/usr/local/bin/gipny-relay
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/gipny-relay

[Install]
WantedBy=multi-user.target
